# deploy/systemd/wechat_artical-feedback-sync.timer.template
[Unit]
Description=wechat_artical feedback sync timer (every 6 hours)
Requires=wechat_artical-feedback-sync.service

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true

[Install]
WantedBy=timers.target
